# scripts/com.nihongo.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nihongo.watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/linh.nguyenthi/Personal/my-nihongo/venv/bin/python</string>
        <string>/Users/linh.nguyenthi/Personal/my-nihongo/processor/watcher.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/linh.nguyenthi/Personal/my-nihongo</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/linh.nguyenthi/Personal/my-nihongo/processor/logs/watcher.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/linh.nguyenthi/Personal/my-nihongo/processor/logs/watcher-error.log</string>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
